# metrics/statsd_5_minute.service
[Unit]
Decription=Run a set of metrics every 5 minutes
After=network.target
[Service]
Type=oneshot
ExecStart=/sbin/metrics_5min

